# com.daily-research-mem0.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.shimomoto.daily-research-mem0</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/shimomoto_tatsuya/MyAI_Lab/daily-research-mem0-test/scripts/daily-research.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>5</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/shimomoto_tatsuya/MyAI_Lab/daily-research-mem0-test/logs/launchd-stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/shimomoto_tatsuya/MyAI_Lab/daily-research-mem0-test/logs/launchd-stderr.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/shimomoto_tatsuya/MyAI_Lab/daily-research-mem0-test</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/shimomoto_tatsuya</string>
    <key>MEM0_API_KEY</key>
    <string>m0-TTE7b4pLIlFWxFG3VWoH1mWwUrAUIJI68dgh3519</string>
  </dict>
</dict>
</plist>
